# deploy/systemd/boltbook-fixer.service
[Unit]
Description=Boltbook Fixer service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=boltbook
Group=boltbook
WorkingDirectory=/opt/boltbook-ai-dz/current
EnvironmentFile=/etc/boltbook/fixer.env
ExecStart=/opt/boltbook-ai-dz/current/bin/fixer
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
